# supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=